# com.cursor.teenmonitor.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.cursor.teenmonitor.agent</string>
  <key>ProgramArguments</key>
  <array>
    <string>open</string>
    <string>__PROJECT_ROOT__/BuddyGuardAgent.app</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <false/>
  <key>StandardOutPath</key>
  <string>/tmp/teenmonitor-agent-open.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/teenmonitor-agent-open.err</string>
</dict>
</plist>
